# gatomalo.service
[Unit]
Description="GatoMalo Server"

[Service]
Environment="ZOHO_AUTH=3978e2966ba80474281580caf6c99e44"
Environment="ZOHO_ORG=667429349"
Environment="ADMIN_USERNAME=Bluetide"
Environment="ADMIN_PASSWORD=megustaresolver"
Type=simple
User=btc
WorkingDirectory=/gatomalo
ExecStart=/gatomalo/venv/bin/python /gatomalo/app.py

[Install]
WantedBy=multi-user.target
